# supervisord.conf
; ============================================================
; supervisord.conf — monolit container
; Runs: nginx + backend + identity-service + lms-service
; ============================================================
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
loglevel=info

[program:nginx]
command=nginx -g "daemon off;"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autorestart=true
priority=10

[program:backend]
command=node server.js
directory=/app
environment=NODE_ENV="production"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autorestart=true
priority=20
startretries=5

[program:identity]
; .NET Identity Service
command=dotnet identity-service.dll
directory=/app/identity
environment=ASPNETCORE_ENVIRONMENT="Production"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autorestart=true
priority=30

[program:lms]
; Java LMS Service
command=java -jar app.jar
directory=/app/lms
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autorestart=true
priority=30
